# munki_pkgs/RolzogCheck/scripts/postinstall
#!/bin/bash
#Script to load my launchdaemon

/bin/launchctl load /Library/LaunchDaemons/com.thoughtworks.rolzogcheck.plist
